# x86_image_wizard/debian/anurad.service
[Unit]
Description=Anura PTY Daemon
After=multi-user.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/bin/anurad

[Install]
WantedBy=multi-user.target
