# x86_image_wizard/alpine/anura-boot
#!/sbin/openrc-run

description="Run Anura Boot Scripts"

start() {
  ebegin "Running Anura Boot Scripts"
  mount / -o remount,rw
  mount -t 9p host9p /root &
  modprobe ne2k-pci
  dhcpcd &
  eend $?
}

stop() {
  ebegin "Stopping Anura Boot Scripts"
  eend $?
}
