# x86_image_wizard/alpine/anura-boot
#!/sbin/openrc-run

description="Run Anura Boot Scripts"

start() {
  rc-service --ifstopped udev start
  rc-service --ifstopped udev-trigger start
  rc-service --ifstopped udev-settle start
  rc-service --ifstopped udev-postmount start
  ebegin "Running Anura Boot Scripts"
  mount / -o remount,rw
  printf "anura" > /etc/hostname
  printf "nameserver 192.168.86.1\n" > /etc/resolv.conf
  rc-service --ifstopped hostname start
  mount -t 9p host9p /root &
  modprobe ne2k-pci
  ifupdown ifup eth0
  ip link set lo up
  eend $?
}

stop() {
  ebegin "Stopping Anura Boot Scripts"
  eend $?
}
